# roles/arch/templates/dock-udev.rules.j2
KERNEL=="card0", ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/{{ user }}/bin/hotplug_monitor"
